# templates/etc/systemd/system/rabbitmq-server.service.d/limits.conf
[Service]
LimitNOFILE = {{ rabbitmq_fd_limits | default(1024) }}
